# infrastructure/docker-containers/maxscale/assets/etc/supervisor/conf.d/maxscale.ini
[program:maxscale]
command=maxscale -d -U maxscale -l stdout
user=maxscale
numprocs=1
stdout_logfile=/var/log/maxscale/stdout_worker.log
stderr_logfile=/var/log/maxscale/stderr_worker.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 60
stopasgroup=true
priority=1000
